# server/supervisord.conf
[program:chrome_call_claude_server]
command=node index.js --port 12346 --log-dir /tmp/chrome/
environment=NODE_ENV=production
directory=/var/www/chrome_call_your_claude_code/server
stopwaitsecs=5
restart=always
stderr_logfile=/tmp/chrome/supervisor-chrome-call-claude-err.log
stdout_logfile=/tmp/chrome/supervisor-chrome-call-claude-out.log
user=root
